# rootfs/etc/s6-overlay/s6-rc.d/salt-api/run
#!/command/execlineb -P

salt-api -ldebug
