# rootfs/etc/s6-overlay/s6-rc.d/generate-certificates/up
#!/command/execlineb -P

salt-run -linfo salt.cmd tls.create_self_signed_cert
